# ops/launchd/com.vojow.reddit-ai-agents-digest.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.vojow.reddit-ai-agents-digest.daily</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/wojciechwieczorek/Sii/reddit-ai-agents-digest/scripts/run_ai_launchd.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/wojciechwieczorek/Sii/reddit-ai-agents-digest</string>

  <key>RunAtLoad</key>
  <false/>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/wojciechwieczorek/Library/Logs/reddit-ai-agents-digest/daily.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/wojciechwieczorek/Library/Logs/reddit-ai-agents-digest/daily.err.log</string>
</dict>
</plist>
